# 95-monitor.rules
KERNEL=="card0", SUBSYSTEM=="drm", ENV{WAYLAND_DISPLAY}="wayland-0", RUN+="/usr/local/bin/SwayMonMan"
